# crontab
Adicionalmente agregamos la tarea al crontab para ejecutarlo cada minuto

* * * * * sh /direccion/del/archivo/checkstatus.sh  >/dev/null 2>&1
